# devconfig/.config/systemd/user/emacs.socket
[Socket]
ListenStream=/home/zzapotheosis/.emacs.d/.socket
#DirectoryMode=0700

[Install]
WantedBy=sockets.target